# packaging/cfengine-nova/debian/cfengine-nova.postinst
#!/bin/sh

#
# Register CFEngine initscript, if not yet.
#
if [ -x /etc/init.d/cfengine3 ]; then
  update-rc.d cfengine3 defaults
fi

#
# Generate a host key
#
if [ ! -f /var/cfengine/ppkeys/localhost.priv ]; then
    /var/cfengine/bin/cf-key >/dev/null || :
fi

exit 0
